# supervisor_mailsystem.ini
[program:message_dispatch]
directory = /www/mailsystem
command = /usr/bin/php artisan queue:work --queue=message-dispatch
autostart = true 
startsecs = 5
autorestart = true
startretries = 3
user = root 
redirect_stderr = true
stdout_logfile_maxbytes = 20MB
stdout_logfile_backups = 20
stdout_logfile = /var/log/supervisor/message_dispatch.log

[program:social_dispatch]
directory = /www/mailsystem
command = /usr/bin/php artisan queue:work --queue=social-dispatch
autostart = true 
startsecs = 5
autorestart = true
startretries = 3
user = root 
redirect_stderr = true
stdout_logfile_maxbytes = 20MB
stdout_logfile_backups = 20
stdout_logfile = /var/log/supervisor/social_dispatch.log

[program:webhook_process]
directory = /www/mailsystem
command = /usr/bin/php artisan queue:work --queue=webhook-process
autostart = true 
startsecs = 5
autorestart = true
startretries = 3
user = root 
redirect_stderr = true
stdout_logfile_maxbytes = 20MB
stdout_logfile_backups = 20
stdout_logfile = /var/log/supervisor/webhook_process.log